# services/systemd/htpc-reset-rclone-mounts.service
[Unit]
Description=Reset HTPC rclone mounts
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
EnvironmentFile=/etc/htpc/htpc.env
ExecStart=/usr/local/bin/htpc-reset-rclone-mounts
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
